# code/cher/docker/etc/supervisor/conf-available/sshd.conf
[program:sshd]
command=/bin/sh -c "mkdir -p /run/sshd; exec /usr/sbin/sshd -p 2222 -D -e"
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/sshd.log
